# etc/supervisor.d/caddy.conf
[program:caddy]
command=caddy -conf /etc/Caddyfile -port 80 %(ENV_CADDY_OPTIONS)s
